# docker/supervisord/api.conf
[program:djangorest]
directory=/var/www/bwd/api
command=python3 manage.py runserver 0.0.0.0:8000
user=root
autostart=true
autorestart=true