# docker/supervisord/client.conf
[program:nuxt]
directory=/var/www/bwd/client
command=node .output/server/index.mjs
user=root
autostart=true
autorestart=true